# plist-templates/self-initiate.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.{BOT_ID}-self-initiate</string>

    <!-- 每 10 分钟 tick，脚本内部随机 30min-24h 自节流（开口时机不可预测） -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>{REPO_DIR}/scripts/self_initiate.py</string>
        <string>{BOT_ID}</string>
        <string>{CHAT_ID}</string>
    </array>

    <key>StartInterval</key>
    <integer>600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>AbandonProcessGroup</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>{REPO_DIR}</string>

    <key>StandardOutPath</key>
    <string>/tmp/{BOT_ID}-self-initiate.stdout</string>
    <key>StandardErrorPath</key>
    <string>/tmp/{BOT_ID}-self-initiate.stderr</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:~/.bun/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
